# dev_server_setup/code/home/anikan/greeter/current/build/upstart/greeter.conf
description "Greeter"

start on runlevel [345]
stop on runlevel [06] or stopped greeter_nginx or stopped greeter_puma