# dev_server_setup/code/home/anikan/greeter/current/build/upstart/greeter_puma.conf
description "Greeter puma"
start on started greeter
stop on stopped greeter

respawn
respawn limit 10 5

exec /bin/su - anikan -c "RACK_ENV=production exec bash -lc '/home/anikan/greeter/current/bin/greeter 9292'"
