# dev_server_setup/code/home/anikan/greeter/current/build/upstart/greeter_nginx.conf
description "Greeter nginx"
start on started greeter
stop on stopped greeter

env CMD="/usr/sbin/nginx -p /home/anikan/greeter/shared/nginx -c /home/anikan/greeter/current/build/nginx/nginx.conf"

expect fork
respawn
respawn limit 10 5

pre-start script
        $CMD -t
        if [ $? -ne 0 ]; then
          exit $?
        fi
end script

exec $CMD

